# etc/hotplug.d/button/50-wps
#
# Copyright (c) 2013 Qualcomm Atheros, Inc..
#
# All Rights Reserved.
# Qualcomm Atheros Confidential and Proprietary.
#

if [ "$ACTION" = "released" -a "$BUTTON" = "wps" ]; then
	echo "" > /dev/console
	echo "WPS PUSH BUTTON EVENT DETECTED, GENERATING SYSEVENT" > /dev/console
        /sbin/sysevent set wps_hw_button pressed
fi
